# systemd/cargo-clean.timer
[Unit]
Description=Clean cargo projects from ~/Git daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
